# docker/crontabs/default
# do daily/weekly/monthly maintenance
# min   hour    day     month   weekday command
* * * *	*  (cd /usr/share/nginx/html; php artisan schedule:run >> /dev/null 2>&1)
